# boxes/privilege_escalation/task2/crontab
SHELL=/bin/sh

*/1 * * * * root touch /home/pupil/it_works.txt
